# docker/supervisord-octane.conf
; ============================================================
; Supervisor untuk FrankenPHP + Octane
; ============================================================
; FrankenPHP web server dijalankan oleh entrypoint.sh langsung
; (bukan supervisor). Supervisor hanya mengelola:
; 1. Queue worker (WhatsApp notif, AI review, password reset)
; 2. Laravel Pulse check
; ============================================================

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:queue-worker]
command=php /app/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/queue.err.log
stdout_logfile=/var/log/supervisor/queue.out.log

[program:pulse]
command=php /app/artisan pulse:check
autostart=true
autorestart=true
user=www-data
stderr_logfile=/var/log/supervisor/pulse.err.log
stdout_logfile=/var/log/supervisor/pulse.out.log